# octopus.conf
[program:octopus]
command = /webapps/octopus/env/bin/gunicorn_start.sh ; Command to start app
user = octopus ; User to run as
stdout_logfile = /webapps/octopus/env/logs/gunicorn_supervisor.log ; Where to write log messages
redirect_stderr = true ; Save stderr in the same log
autorestart=true